# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:pentaho]
command=/opt/app-root/src/biserver-ce/start-pentaho.sh
autorestart=false
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
priority=1
exitcodes=0
startsecs=0
